# etc/profile.d/10linuxlogo.sh
[ -x /usr/bin/linuxlogo ] && /usr/bin/linuxlogo
