# etc/profile.d/99bash_tricks.sh
# Auto change directory.
shopt -s autocd

# Mistyped paths auto-correction.
shopt -s cdspell

# Case insensitive autocompletion.
shopt -s nocaseglob


